# mc1.17.1-restart.timer
[Unit]
Description=senkamc Server Restart Schedule Task

[Timer]
OnCalendar=*-*-* 04:59:30

[Install]
WantedBy=multi-user.target
